# deployment/files/supervisor/postgresql.conf
[program:postgresql]
;user = postgres
command = /usr/bin/postmaster -D "/var/lib/postgres/data"
process_name = %(supervisor_process_base)s
stopsignal = INT
autostart = true
autorestart = true
