# config/etc/cron.d/mturk
APP_PATH=/var/www/mturk
MAILTO="michal.klujszo@10clouds.com"

#0       *       *       *       *       root    cd $APP_PATH && python manage.py crawler 1
15      *       *       *       *       root    cd $APP_PATH && python manage.py db_refresh_mviews
25      1       *       *       *       root    cd $APP_PATH && python manage.py db_calculate_daily_stats

