# config/etc/cron.daily/psqlmntn
#!/bin/sh
su -c 'vacuumdb --all --full --analyze' postgres
su -c 'reindexdb --all' postgres